# alpine/5/rootfs/etc/services.d/supervisord/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid root
docker-entrypoint # This runs the original entrypoint logic
supervisord -c /etc/supervisord.conf